# devopsfetch.service
[Unit]
Description=DevOpsFetch Monitoring Service
After=network.target

[Service]
ExecStart=/path/to/devopsfetch/devopsfetch.sh -l
Restart=always
User=root

[Install]
WantedBy=multi-user.target

